# system/backends/void/runit/alpenglow-netd/run
#!/bin/sh
set -eu
exec /usr/local/bin/alpenglow-netd
