# system/alpine/openrc/alpenglow-netd
#!/sbin/openrc-run

description="Alpenglow network state daemon"

command="/usr/local/bin/alpenglow-netd"
command_background="yes"
pidfile="/run/${RC_SVCNAME}.pid"
export ALPENGLOW_NETD_STATE_JSON="${ALPENGLOW_NETD_STATE_JSON:-/run/alpenglow/netd/interfaces.json}"
export ALPENGLOW_NETD_RUNTIME_ENV="${ALPENGLOW_NETD_RUNTIME_ENV:-/run/alpenglow/netd/runtime-state.env}"
export ALPENGLOW_NETD_INTERVAL_SECS="${ALPENGLOW_NETD_INTERVAL_SECS:-5}"
export ALPENGLOW_RUNTIME_EVENTS_FILE="${ALPENGLOW_RUNTIME_EVENTS_FILE:-/run/alpenglow/runtime-events.log}"

now_unix_ms() {
    printf '%s000\n' "$(date +%s)"
}

record_runtime_event() {
    kind="$1"
    source="$2"
    detail="${3:-}"
    dir="$(dirname "${ALPENGLOW_RUNTIME_EVENTS_FILE}")"
    mkdir -p "${dir}"
    printf '%s\t%s\t%s\t%s\n' "$(now_unix_ms)" "${kind}" "${source}" "${detail}" >>"${ALPENGLOW_RUNTIME_EVENTS_FILE}" 2>/dev/null || true
}

start_pre() {
    record_runtime_event network_daemon_start openrc alpenglow-netd
}

start_post() {
    record_runtime_event network_daemon_started openrc alpenglow-netd
}

depend() {
    need localmount networking alpenglow-kernel-policy
    before sold alpenglow-session
}
