# system/alpine/openrc/alpenglow-pressure
#!/sbin/openrc-run

description="Alpenglow pressure policy"

command="/usr/local/bin/apply-pressure-policy.sh"
command_background="yes"
pidfile="/run/alpenglow/alpenglow-pressure.pid"

depend() {
    need localmount alpenglow-kernel-policy alpenglow-zram
    before sold alpenglow-session
}
