# system/alpine/openrc/alpenglow-zram
#!/sbin/openrc-run

description="Alpenglow zram swap"

command="/usr/local/bin/apply-zram-policy.sh"
command_background="no"

depend() {
    need localmount alpenglow-kernel-policy
    before sold alpenglow-session
}
