# system/alpine/openrc/alpenglow-session
#!/sbin/openrc-run

description="Alpenglow session launcher"

command="/usr/local/bin/alpenglow-session-start"

depend() {
    need sold seatd alpenglow-kernel-policy alpenglow-zram
}

respawn=YES
respawn_delay=5
