# system/alpine/openrc/sold
#!/sbin/openrc-run

description="Alpenglow local server"

command="/usr/local/bin/sold"
command_background="yes"
pidfile="/run/${RC_SVCNAME}.pid"
export ALPENGLOW_BUNDLE_DIR="/usr/local/share/alpenglow/bundle"
export SOLD_BIND="127.0.0.1:8080"
export ALPENGLOW_SYSTEM_CONFIG="/etc/alpenglow/system.json"
export ALPENGLOW_PLUGIN_STATE_FILE="/var/lib/alpenglow/system/plugin-state.json"
export ALPENGLOW_PLUGIN_INSTALL_STATE_FILE="/var/lib/alpenglow/system/plugin-installs.json"
export ALPENGLOW_PLUGIN_MANIFEST_DIR="/etc/alpenglow/plugins"
export ALPENGLOW_UPDATE_POLICY_FILE="/etc/alpenglow/update-policy.json"
export ALPENGLOW_UPDATE_STATE_FILE="/var/lib/alpenglow/system/update-state.json"
export ALPENGLOW_PACKAGE_MANAGER_FILE="/etc/alpenglow/package-manager.json"
export ALPENGLOW_RUNTIME_STATE_ENV="${ALPENGLOW_RUNTIME_STATE_ENV:-/run/alpenglow/runtime-state.env}"
export ALPENGLOW_RUNTIME_EVENTS_FILE="${ALPENGLOW_RUNTIME_EVENTS_FILE:-/run/alpenglow/runtime-events.log}"

now_unix_ms() {
    printf '%s000\n' "$(date +%s)"
}

record_runtime_state() {
    key="$1"
    value="$2"
    dir="$(dirname "${ALPENGLOW_RUNTIME_STATE_ENV}")"
    tmp="${ALPENGLOW_RUNTIME_STATE_ENV}.$$"
    mkdir -p "${dir}"
    if [ -f "${ALPENGLOW_RUNTIME_STATE_ENV}" ]; then
        grep -v "^${key}=" "${ALPENGLOW_RUNTIME_STATE_ENV}" >"${tmp}" || true
    else
        : >"${tmp}"
    fi
    printf '%s=%s\n' "${key}" "${value}" >>"${tmp}"
    mv "${tmp}" "${ALPENGLOW_RUNTIME_STATE_ENV}"
}

record_runtime_event() {
    kind="$1"
    source="$2"
    detail="${3:-}"
    dir="$(dirname "${ALPENGLOW_RUNTIME_EVENTS_FILE}")"
    mkdir -p "${dir}"
    printf '%s\t%s\t%s\t%s\n' "$(now_unix_ms)" "${kind}" "${source}" "${detail}" >>"${ALPENGLOW_RUNTIME_EVENTS_FILE}" 2>/dev/null || true
}

attach_to_cgroup() {
    group="$1"
    pid="$2"
    if [ -x /usr/local/bin/alpenglow-kernelctl ]; then
        /usr/local/bin/alpenglow-kernelctl attach --group "${group}" --pid "${pid}" >/dev/null 2>&1 || true
        return
    fi
    path="/sys/fs/cgroup/alpenglow/${group}"
    if [ -f /sys/fs/cgroup/cgroup.controllers ] && mkdir -p "${path}" 2>/dev/null; then
        printf '%s\n' "${pid}" >"${path}/cgroup.procs" 2>/dev/null || true
    fi
}

start_pre() {
    record_runtime_state ALPENGLOW_SOLD_START_UNIX_MS "$(now_unix_ms)"
    record_runtime_event sold_openrc_start openrc sold
}

start_post() {
    if [ -f "${pidfile}" ]; then
        attach_to_cgroup system "$(cat "${pidfile}")"
    fi
    record_runtime_event sold_openrc_started openrc sold
}

depend() {
    need localmount networking alpenglow-kernel-policy alpenglow-netd alpenglow-zram
    after seatd
}
