# system/alpine/rootfs-overlay/etc/inittab
# /etc/inittab for Alpenglow v0

::sysinit:/sbin/openrc sysinit
::sysinit:/sbin/openrc boot
::wait:/sbin/openrc default
::shutdown:/sbin/openrc shutdown
